# deploy/launchd/com.nfc-scan.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.nfc-scan.backend</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__APP_DIR__/deploy/launchd/run-backend.sh</string>
  </array>
  <key>WorkingDirectory</key><string>__APP_DIR__</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>__APP_DIR__/backend.launchd.log</string>
  <key>StandardErrorPath</key><string>__APP_DIR__/backend.launchd.log</string>
</dict>
</plist>
